# com.chaitin.quarterly-estimate-push.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.chaitin.quarterly-estimate-push</string>
    <key>ProgramArguments</key>
    <array>
        <string>REPLACE_WITH_SKILL_DIR/scripts/run.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key><integer>8</integer>
        <key>Minute</key><integer>40</integer>
    </dict>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>REPLACE_WITH_HOME</string>
    </dict>
    <key>StandardOutPath</key>
    <string>REPLACE_WITH_SKILL_DIR/logs/stdout.log</string>
    <key>StandardErrorPath</key>
    <string>REPLACE_WITH_SKILL_DIR/logs/stderr.log</string>
</dict>
</plist>
